# docs/systemd/tg-sms-bot.service
[Unit]
Description=TG SMS Bot Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/tg-sms-suite
Environment=VENV_DIR=/opt/tg-sms-suite/.venv
Environment=LOG_DIR=/var/log/tg-sms-suite
ExecStart=/opt/tg-sms-suite/scripts/run_bot.sh
Restart=on-failure
RestartSec=5
User=tgsms
Group=tgsms

[Install]
WantedBy=multi-user.target
